# env/application.ini
[program:httpd]
command = /usr/libexec/s2i/run
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
;user=default

;[program:php-fpm]
;command = /usr/sbin/php-fpm
;stdout_logfile=/dev/stdout
;stdout_logfile_maxbytes=0
;stderr_logfile=/dev/stderr
;stderr_logfile_maxbytes=0
;user=default
[program:postgres]
command = postgresql-15-setup --initdb